# supervisor/conf.d/app.conf
[program:test]
command=php artisan queue:work

process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
stdout_logfile=/app/storage/logs/test.log
